# _docker/app/supervisor_my.conf
[program:my_process]
command=/usr/local/bin/php console consume
directory=/var/www
stdout_logfile=/var/log/supervisor/my_process.log 
redirect_stderr=true